# deploy/ticket-monitor.service
[Unit]
Description=Ticketmaster Ticket Availability Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ticket-monitor
Group=ticket-monitor
WorkingDirectory=/opt/ticket-monitor

# systemd creates /run/ticket-monitor before the service starts, owned by
# User:Group, and wipes it on stop. Avoids the "non-existent path" error
# that ReadWritePaths= would otherwise raise on first boot.
RuntimeDirectory=ticket-monitor
RuntimeDirectoryMode=0750

# "+" prefix: run as root, bypassing User= and the sandbox. fetch-secrets.sh
# needs root to write an env file owned by root:ticket-monitor mode 640.
ExecStartPre=+/opt/ticket-monitor/deploy/fetch-secrets.sh

# Leading "-": tolerate a missing file. On first start the file doesn't
# exist yet; ExecStartPre creates it, and ExecStart re-reads it.
EnvironmentFile=-/run/ticket-monitor/env

ExecStart=/opt/ticket-monitor/.venv/bin/python3 /opt/ticket-monitor/monitor.py

Restart=on-failure
RestartSec=15s
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
